# roles/arch/files/glances.service
[Unit]
Description=Glances Server

[Service]
ExecStart=/etc/glances/bin/glances --config /etc/glances/glances.conf -s

[Install]
WantedBy=multi-user.target
