# ansible/roles/root/files/glances.service
[Unit]
Description=Glances Server

[Service]
ExecStart=/root/.local/bin/glances -w

[Install]
WantedBy=multi-user.target
